# telegram-project/fewijhca3fih4bot.service
[Unit]
Description=Telegram hook listener for @fewijhca3fih4bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/fewijhca3fih4bot
EnvironmentFile=/opt/fewijhca3fih4bot/.env
ExecStart=/opt/fewijhca3fih4bot/.venv/bin/python /opt/fewijhca3fih4bot/reply_ok_bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
